# debian/postinst
#!/bin/bash
update-rc.d jrs-daemon defaults 99
